# contrib/dist/el7/broker.service
[Unit]
Description=The Choria Orchestrator Middleware Broker
After=network.target

[Service]
StandardOutput=syslog
StandardError=syslog
User=root
Group=root
ExecStart={{bindir}}/{{pkgname}} broker --config={{etcdir}}/broker.conf

[Install]
WantedBy=multi-user.target
