# packager/templates/debian/generic/prerm
#!/bin/sh

action="$1"

if [ "$action" = remove ]; then
	systemctl --no-reload disable choria-server.service choria-broker.service
	systemctl stop choria-server.service choria-broker.service
fi
